# scripts/com.grillmaster.tunnel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.grillmaster.tunnel</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>PORT=8766; DATA_DIR="/Users/admin/Documents/GitHub/grillmaster-command-center/data"; URL_FILE="$DATA_DIR/remote-url.txt"; LOG_FILE="$DATA_DIR/localhost-run.log"; mkdir -p "$DATA_DIR"; : > "$LOG_FILE"; rm -f "$URL_FILE"; exec /usr/bin/ssh -oExitOnForwardFailure=yes -oServerAliveInterval=30 -oServerAliveCountMax=3 -oStrictHostKeyChecking=no -oLogLevel=ERROR -R 80:localhost:${PORT} nokey@localhost.run 2>&1 | /usr/bin/tee -a "$LOG_FILE" | /usr/bin/awk -v urlfile="$URL_FILE" '"'"'/tunneled with tls termination, https:\/\// { if (match($0, /https:\/\/[A-Za-z0-9.-]+/)) { print substr($0, RSTART, RLENGTH) > urlfile; close(urlfile) } } { print; fflush() }'"'"'</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/admin</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
